# devops/local-kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
